# systemd/oceano-player.service
[Unit]
Description=Oceano Player (AirPlay receiver supervisor)
After=network.target sound.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/opt/oceano-player
ExecStart=/opt/oceano-player/bin/oceano-player -config /opt/oceano-player/config.yaml
Restart=always
RestartSec=2

# Hardening (keep modest; audio daemons sometimes need capabilities)
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target

