# packaging/systemd/oceano-state-manager.service
[Unit]
Description=Oceano State Manager (unified playback state + ACRCloud recognition)
After=network-online.target shairport-sync.service oceano-source-detector.service
Wants=shairport-sync.service network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/oceano-state-manager \
  --metadata-pipe "/tmp/shairport-sync-metadata" \
  --source-file "/tmp/oceano-source.json" \
  --output "/tmp/oceano-state.json" \
  --artwork-dir "/var/lib/oceano/artwork" \
  --pcm-socket "/tmp/oceano-pcm.sock" \
  --vu-socket "/tmp/oceano-vu.sock" \
  --recognizer-capture-duration "7s" \
  --recognizer-max-interval "5m0s" \
  --idle-delay "10s" \
  --library-db "/var/lib/oceano/library.db"
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
