# packaging/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Runtime directories (nfpm creates them, but ensure correct perms on upgrade)
        install -d -m 755 /etc/oceano
        install -d -m 755 /var/lib/oceano
        install -d -m 755 /var/lib/oceano/artwork
        install -d -m 755 /usr/local/lib/oceano

        # Write default config only on first install — never overwrite on upgrade
        if [ ! -f /etc/oceano/config.json ]; then
            cp /etc/oceano/config.json.default /etc/oceano/config.json
            chmod 0644 /etc/oceano/config.json
        fi

        # Python virtualenv for Broadlink IR bridge (non-fatal if python3 unavailable)
        if command -v python3 >/dev/null 2>&1; then
            if [ ! -d /opt/oceano-venv ]; then
                python3 -m venv /opt/oceano-venv
            fi
            if ! /opt/oceano-venv/bin/python3 -c "import broadlink" 2>/dev/null; then
                /opt/oceano-venv/bin/pip install --quiet broadlink || true
            fi
        fi

        # Enable and start services
        systemctl --no-ask-password daemon-reload
        systemctl --no-ask-password enable oceano-source-detector.service
        systemctl --no-ask-password enable oceano-state-manager.service
        systemctl --no-ask-password enable oceano-web.service
        systemctl --no-ask-password restart oceano-source-detector.service || true
        systemctl --no-ask-password restart oceano-state-manager.service   || true
        systemctl --no-ask-password restart oceano-web.service             || true

        IP=$(hostname -I 2>/dev/null | awk '{print $1}') || IP="<pi-ip>"
        echo ""
        echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
        echo "  Installation Complete"
        echo "━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━━"
        echo ""
        echo "Open http://${IP}:8080 to configure audio devices and ACRCloud."
        echo ""
        echo "For AirPlay, Bluetooth, ALSA devices, and optional HDMI/DSI kiosk, run:"
        echo "  sudo oceano-setup"
        echo ""
        echo "See the project README for capture gain calibration and advanced display options."
        echo ""
        echo "Monitor logs:"
        echo "  journalctl -u oceano-web.service -f"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
